# .config/systemd/user/apps.target
[Unit]
Description=Desktop Apps

[Install]
WantedBy=wm.target
