# .config/systemd/user/borg.timer
[Unit]
Description=Backup timer

[Timer]
OnActiveSec=6h

[Install]
WantedBy=default.target
