# .config/systemd/user/i3@.socket
[Unit]
Description=Socket for i3 at display %i

[Socket]
SocketUser=tim
ListenStream=/run/user/1000/sway/ipc.%i.sock
Service=sway@%i.service
