# .config/systemd/user/notmuch.timer
[Unit]
Description=Index Maildir

[Timer]
OnUnitActiveSec=1h
Unit=notmuch.service

[Install]
WantedBy=default.target
