# .config/systemd/user/offlineimap.timer
[Unit]
Description=Sync Maildir

[Timer]
OnUnitActiveSec=1h
Unit=offlineimap.service

[Install]
WantedBy=default.target
