# .config/systemd/user/redshift.service
[Unit]
Description=Redshift
After=graphical.target

[Service]
ExecStart=/usr/bin/redshift -l geoclue2

[Install]
WantedBy=graphical.target