# .config/systemd/user/wm.target
[Unit]
Description=Window manager target 

[Install]
WantedBy=graphical.target
