# .config/systemd/system/mnt-backup.automount
[Unit]
Description=Mounts backup drive

[Automount]
Where=/mnt/backup

[Install]
WantedBy=multi=user.target
